# systemd/mpdScheduler.service
[Unit]
Description=mpdScheduler

[Service]
ExecStart=/usr/bin/mpdScheduler

[Install]
WantedBy=multi-user.target
