# systemd/mpdScheduler.service
[Unit]
Description=mpdScheduler

[Service]
ExecStart=WHICH_MPDSCHEDULER

[Install]
WantedBy=multi-user.target
